# squid/ut1-to-sqlite.timer
[Unit]
Description=Daily update UT1 database

[Timer]
OnCalendar=*-*-* 03:40:00
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
